# cron/cron.d/gcr-cleanup
SHELL=/bin/sh

7 * * * * gcr-cleanup nice -19 /usr/local/bin/gcr-cleanup.sh
